# resources/cron/crontab
* * * * * loxberry REPLACELBPBINDIR/cron_runner
